# supervisord-etcd-service-registry.conf
[program:etcd-service-registry]
command=/start-etcd-service-registry.sh
numprocs=1
autostart=true
autorestart=true
